# infra/files/unicoach-rest-server.service
[Unit]
Description=unicoach REST API server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=unicoach
EnvironmentFile=/etc/unicoach/env
ExecStart=/opt/unicoach/current/rest-server/bin/rest-server
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
